# etc/local.d/99_led_config.start
echo mmc1 > /sys/class/leds/ack/trigger
